# deployment/com.business-stories.orchestrator.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!--
        macOS LaunchAgent: fires run_orchestrator.sh every 5 minutes.

        Each firing executes ONE pipeline stage and exits. The pipeline's
        file-lock (state/locks/orchestrator.lock) prevents concurrent
        runs: if S09 is mid-render, the next firing logs "lock contention"
        and exits cleanly. So scheduling at a high frequency is harmless
        — most firings will be no-ops while a slow stage runs.

        Install:
            cp deployment/com.business-stories.orchestrator.plist \
               ~/Library/LaunchAgents/
            launchctl load -w ~/Library/LaunchAgents/com.business-stories.orchestrator.plist

        Check status:
            launchctl list | grep business-stories

        Stop temporarily:
            launchctl unload ~/Library/LaunchAgents/com.business-stories.orchestrator.plist

        Remove permanently:
            launchctl unload -w ~/Library/LaunchAgents/com.business-stories.orchestrator.plist
            rm ~/Library/LaunchAgents/com.business-stories.orchestrator.plist

        IMPORTANT — replace the two ABSOLUTE PATHS below ("/Users/cantemir/...")
        with the actual paths on your machine before copying into
        ~/Library/LaunchAgents/. launchd will not honour ${...} substitution.
    -->

    <key>Label</key>
    <string>com.business-stories.orchestrator</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/cantemir/Projects/business_success_stories/run_orchestrator.sh</string>
    </array>

    <!-- 300s = 5 min. Drop to 60 if you want minute-level cadence
         while debugging; raise to 600 (10 min) once the pipeline is
         running steadily and you don't need rapid feedback. -->
    <key>StartInterval</key>
    <integer>300</integer>

    <!-- Fire once immediately on launchctl load so you don't wait 5
         min for the first run. -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Working directory matters: pipeline/__init__.py looks for
         .env at the CWD, and config.yaml is resolved relative to
         the project root anyway. -->
    <key>WorkingDirectory</key>
    <string>/Users/cantemir/Projects/business_success_stories</string>

    <!-- launchd's PATH is minimal by default — set what we need.
         /opt/homebrew/bin for Apple Silicon brew (ffmpeg, ffprobe,
         flux, python3), /usr/local/bin for Intel brew, /usr/bin for
         /usr/bin/env. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <!-- LC_ALL forces utf-8 so subprocess output from FLUX/ffmpeg
             doesn't blow up on non-ASCII bytes. -->
        <key>LC_ALL</key>
        <string>en_US.UTF-8</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
    </dict>

    <!-- launchd appends to these files. The orchestrator does its
         own structured logging into logs/orch.YYYY-MM-DD.log via
         pipeline/hermes_orchestrator.py — these capture only
         launchd-level events (process spawn, exit codes, accidental
         stderr) so we can tell when the agent itself misbehaved. -->
    <key>StandardOutPath</key>
    <string>/Users/cantemir/Projects/business_success_stories/logs/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/cantemir/Projects/business_success_stories/logs/launchd.err.log</string>

    <!-- Don't auto-restart on exit — a stage finishing IS the
         expected exit, not a crash. StartInterval will re-fire on
         schedule anyway. -->
    <key>KeepAlive</key>
    <false/>

    <!-- Run as the logged-in user, not root. (LaunchAgent in
         ~/Library is already user-scoped, but explicit is better.) -->
    <key>ProcessType</key>
    <string>Background</string>

    <!-- Lower CPU + I/O priority so FLUX renders don't interfere
         with foreground work if you're using the same Mac
         interactively. -->
    <key>Nice</key>
    <integer>5</integer>
</dict>
</plist>
